# crontab
0 * * * * python3 /app/sync.py >> /app/log.txt